# nginx/etc/init/nginx.conf
description 'nginx webserver'

start on startup
stop  on shutdown

respawn
expect fork
exec {{install_dir}}/sbin/nginx
